# example.socket
[Unit]

[Socket]
ListenStream=5432

[Install]
WantedBy=default.target
